# lib/systemd/system/i2p-torbrowser-sockets-workaround_autogen__var_run_i2p-torbrowser-sockets-workaround_127.0.0.1_4447.sock.service
## Copyright (C) 2018 
## See the file COPYING for copying conditions.

## This file was auto generated by ./usr/lib/i2p-torbrowser-sockets-workaround/systemd-unit-files-generator at package build time.

[Unit]
Description=redirect /var/run/i2p-torbrowser-sockets-workaround/127.0.0.1_4447.sock to Whonix-Gateway port 4447
Documentation=https://github.com/eyedeekay/i2p-torbrowser-sockets-workaround
ConditionPathExists=!/var/run/qubes/this-is-templatevm
After=i2p-torbrowser-sockets-workaround_autogen__var_run_i2p-torbrowser-sockets-workaround_127.0.0.1_4447.sock.socket
Requires=i2p-torbrowser-sockets-workaround_autogen__var_run_i2p-torbrowser-sockets-workaround_127.0.0.1_4447.sock.socket

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd 10.152.152.10:4447
PrivateTmp=yes
