# lib/systemd/system/i2p-torbrowser-sockets-workaround_autogen__var_run_i2p_http.service
## Copyright (C) 2018 
## See the file COPYING for copying conditions.

## This file was auto generated by ./usr/lib/i2p-torbrowser-sockets-workaround/systemd-unit-files-generator at package build time.

[Unit]
Description=redirect /var/run/i2p/http to Whonix-Gateway port 4444
Documentation=https://github.com/eyedeekay/i2p-torbrowser-sockets-workaround
ConditionPathExists=!/var/run/qubes/this-is-templatevm
After=i2p-torbrowser-sockets-workaround_autogen__var_run_i2p_http.socket
Requires=i2p-torbrowser-sockets-workaround_autogen__var_run_i2p_http.socket

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd 10.152.152.11:4444
PrivateTmp=yes
