# lib/systemd/system/i2p@default.service.d/51_i2p_torbrowser_sockets_workaround.conf
## Copyright (C) 2017 
## See the file COPYING for copying conditions.

[Unit]
Description=i2p-torbrowser-sockets-workaround
Documentation=https://github.com/eyedeekay/i2p-torbrowser-sockets-workaround

[Service]
Type=simple
ExecStartPre=
ExecReload=
AppArmorProfile=
KillSignal=SIGTERM
