# lib/systemd/system/i2p.service.d/50_i2p_torbrowser_sockets_workaround.conf
## Copyright (C) 2012 
## See the file COPYING for copying conditions.

[Service]
## Make 'sudo service tor status' exit '0' for better compatibility.
RemainAfterExit=yes

[Unit]
## Qubes-Whonix:
## Clear 'ConditionPathExists=/var/run/qubes-service/whonix-gateway' set by
## the qubes-whonix package, which is useful on the gateway but not on the
## workstation.
ConditionPathExists=
